# templates/cmppatcher-rewriter.service.tmpl
[Unit]
Description=cmppatcher FMA Rewriter Daemon (CMP 170HX FFMA bypass)
Documentation=https://github.com/amoghmunikote/170th-street
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 __INSTALL_DIR__/rewriter_daemon.py \
    --socket /var/run/cmppatcher-rewriter.sock \
    --cache-dir __INSTALL_DIR__/cache
Restart=on-failure
RestartSec=5
RuntimeDirectory=cmppatcher
RuntimeDirectoryMode=0755
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cmppatcher-rewriter

[Install]
WantedBy=multi-user.target
